# roles/tls-certificate/templates/etc/systemd/system/acme-sh.timer.j2
[Unit]
Description=acme.sh Timer

[Timer]
OnCalendar=Sun 01:17:00
Unit=acme-sh.service

[Install]
WantedBy=timers.target